# templates/systemd/powerdns_admin.socket.j2
[Unit]
Description={{ pdnsadmin_name }} socket

[Socket]
ListenStream=/run/powerdns-admin/socket
SocketUser={{ pdnsadmin_user }}
SocketGroup={{ pdnsadmin_group }}

[Install]
WantedBy=sockets.target
