# templates/systemd/powerdns-admin.socket.j2
[Unit]
Description={{ pdnsadmin_name }} socket
After=multi-user.target

[Socket]
ListenStream=/run/powerdns-admin/socket

[Install]
WantedBy=sockets.target
